# local_monitoring/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:monitoring]
command=/usr/bin/python /var/www/run.py
autorestart=true
